# playground/scripts/crontab
*/5 * * * * /root/longrun/checkRunningScript.sh >> /root/longrun/log 2>&1
